# ddd.conf
[program:ddd]
command=python3 /home/curator2/ddd/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/ddd/test.err.log
stdout_logfile=/home/curator2/ddd/test.out.log